# deploy/nemulai-agent.service
# NemulAI GPU Energy Monitoring Agent — systemd service unit
#
# Install:
#   sudo cp nemulai-agent.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nemulai-agent
#
# Or use the installer:  nemulai service install
# Or:  curl -sSL https://get.nemulai.com | bash

[Unit]
Description=NemulAI GPU Energy Monitoring Agent
Documentation=https://nemulai.com/docs/agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=nemulai
Group=nemulai

# Read API key and settings from the env file.
# Sensitive values (API key) live here with mode 600.
EnvironmentFile=/etc/nemulai/agent.env

# DATA_DIR and LOG_DIR are set here so the agent writes outside the install dir.
# Override in agent.env if you prefer different paths.
Environment=DATA_DIR=/var/lib/nemulai
Environment=LOG_DIR=/var/log/nemulai

ExecStart=/usr/local/bin/nemulai

# Restart on any non-zero exit, including OOM kills.
Restart=on-failure
RestartSec=10s

# Give the agent time to flush the WAL on shutdown before SIGKILL.
TimeoutStopSec=30s

# Watchdog: systemd kills + restarts the service if it hangs for >120s.
# The agent must call sd_notify(WATCHDOG=1) if WatchdogSec is active.
# Disabled by default; enable once watchdog notify support is added.
# WatchdogSec=120s

# ── Security hardening ────────────────────────────────────────────────────────

# Prevent the agent from gaining new privileges via setuid/setgid.
NoNewPrivileges=true

# Mount /usr, /boot, /efi read-only (agent only reads its own data dir).
ProtectSystem=strict

# Hide /home and /root from the agent process.
ProtectHome=true

# Give the process a private /tmp that is invisible to other users.
PrivateTmp=true

# Allow writes only to the data and log directories.
ReadWritePaths=/var/lib/nemulai /var/log/nemulai

# The agent reads /proc/<pid>/environ and /proc/<pid>/cmdline to attribute jobs.
# /proc is always accessible regardless of ProtectSystem.

# Restrict dangerous kernel calls. "system-service" is a safe broad filter.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Drop ambient capabilities the agent does not need.
CapabilityBoundingSet=

# ── Resource limits ───────────────────────────────────────────────────────────

# Resident memory cap. Increase if running many GPUs with large WAL replay.
MemoryMax=256M

# Enough file descriptors for WAL file + network sockets + Prometheus.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
